# doc/pingu-checks.service
[Unit]
Description=Pingu Check Runner
After=network.target

[Service]
Type=oneshot
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/pingu
ExecStart=/home/exedev/pingu/.venv/bin/python /home/exedev/pingu/manage.py run_checks
