# doc/pingu-checks.timer
[Unit]
Description=Run Pingu checks every minute

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
